# roles/camilladsp/templates/camilladsp.service.j2
[Unit]
Description=CamillaDSP Audio Processing Service
After=network.target sound.target

[Service]
Type=simple
User={{ camilladsp_user }}
Group={{ camilladsp_group }}
ExecStart={{ camilladsp_install_dir }}/bin/camilladsp -p1234 -w -s {{ camilladsp_install_dir }}/statefile.yml
WorkingDirectory={{ camilladsp_install_dir }}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target